# other/current/home-skel/.config/autostart/copyq.desktop
[Desktop Entry]
Exec=copyq
Name=copyq
Type=Application
Version=1.0
