# other/current/home-skel/.config/autostart/i3-appmenu-service.desktop
[Desktop Entry]
Exec=i3-appmenu-service &
Name=i3-appmenu-service
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
