# other/current/home-skel/.config/autostart/kbi keyboard indicator.desktop
[Desktop Entry]
Exec=kbi &
Name=kbi keyboard indicator
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
