# other/current/home-skel/.config/autostart/konsolekill.desktop
[Desktop Entry]
Exec=pkill konsole
Name=killall-konsole
Type=Application
Version=1.0
