# other/current/home-skel/.config/autostart/krunner.desktop
[Desktop Entry]
Exec=krunner -d &
Name=krunner
Type=Application
Version=1.0
