# other/current/home-skel/.config/autostart/motrix.desktop
[Desktop Entry]
Exec=motrix &
Name=motrix
Type=Application
Version=1.0
