# other/current/home-skel/.config/autostart/qmenu.desktop
[Desktop Entry]
Exec=qmenu_registrar &
Name=qmenu
Type=Application
Version=1.0
